# deploy/systemd/salomao-linx-auto-sync@.service
[Unit]
Description=Salomao Linx auto sync (%i)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=salomao
Group=salomao
WorkingDirectory=/srv/salomao/%i/app/backend
ExecStart=/srv/salomao/%i/app/backend/.venv/bin/python -m app.jobs.linx_auto_sync

